# conf/crontab
*/15 * * * * /usr/local/bin/dnstun "UPDATE" >/dev/null 2>&1
