# conf/rc.local
#!/bin/sh -e
/usr/local/bin/dnstun-init &
